# meta-quanta/meta-olympus-nuvoton/recipes-olympus-nuvoton/loadmcu/files/mcu-update.service
[Unit]
Description=NPCM7xx MCU F/W update service

[Service]
ExecStart=/usr/bin/mcu-update.sh
Type=oneshot
